# supervisord.conf
[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

[program-default]
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# attestation server
[program:attestation-server]
command=/app/attestation-server --ip-addr 127.0.0.1:1300 --pub-key /app/id.pub

# attestation server proxy
[program:attestation-proxy]
command=/app/vsock-to-ip --vsock-addr 88:1300 --ip-addr 127.0.0.1:1300

# transparent proxy component inside enclave
[program:ip-to-vsock-transparent]
command=/app/ip-to-vsock-transparent --vsock-addr 3:1200 --ip-addr 127.0.0.1:1200

# # DNS-over-HTTPS provider
# [program:dnsproxy]
# command=/app/dnsproxy -u https://1.1.1.1/dns-query

# # docker
# [program:docker]
# command=/bin/dockerd

# following trader_abci_0 container logs
[program:trader_abci_0-container-logs]
command=docker logs trader_abci_0 --follow
